# root/etc/rc.local
#!/bin/sh -e

# run tor firewall
/usr/local/bin/tor-firewall.sh

exit 0